# sa/fswatch-kepler.conf
[program:fswatch-kepler]
command=fswatch-kepler.sh
directory=/Users/piperck/Desktop/tools/sa
