# sa/fswatch-xcf.conf
[program:fswatch-xcf]
command=./fswatch-xiachufang.sh
directory=/Users/piperck/Desktop/tools/sa
user=zouzhinandeMacBook-Pro
